# .cursor/deployment/server-workspace/openclaw-gateway.service
# ПРИМЕЧАНИЕ: основной [Service] блок управляется OpenClaw автоматически при обновлении.
# Этот файл — СПРАВОЧНАЯ КОПИЯ структуры. Реальные API ключи хранятся в drop-in файлах:
#   ~/.config/systemd/user/openclaw-gateway.service.d/
# Drop-in файлы НЕ перезаписываются при обновлении OpenClaw.

[Unit]
Description=OpenClaw Gateway
After=network-online.target
Wants=network-online.target

[Service]
# ExecStart управляется openclaw install/update — не редактировать вручную
Restart=always
RestartSec=5s
WorkingDirectory=/home/openclaw/.openclaw
StandardOutput=journal
StandardError=journal
MemoryMax=2G
